# apps/plasma/autostart/hyperhdr.desktop
[Desktop Entry]
Exec=hyperhdr
Icon=
Name=hyperhdr
Path=
Terminal=False
Type=Application
