# apps/plasma/autostart/hyperhdr.desktop
[Desktop Entry]
Exec=hyperhdr
Icon=
Name=HyperHDR
Path=
Terminal=false
Type=Application
