# apps/plasma/autostart/hyperhdr.desktop
[Desktop Entry]
Exec="/home/eymeric/.config/autostart/hyperhdr.sh"
Icon=
Name=hyperhdr
Path=
Terminal=False
Type=Application
